# .config/crontab
0 * * * * ~/.local/scripts/update-weather.sh 
30 * * * * sudo xbps-install -S
0 20 * * * ~/.local/scripts/backup.sh
