# .config/crontab
0 * * * * ~/scripts/update-weather.sh 
30 * * * * sudo xbps-install -S
0 20 * * * ~/scripts/backup.sh
